# deploy/systemd/robocar-cam-hub.service
[Unit]
Description=Robocar camera hub (OAK-D owner, publishes frames to /dev/shm)

[Service]
Type=simple
User=robocar
# Défaut prod = clone de référence sur main. sync-services.sh réécrit ce chemin
# vers le repo d'où on lance le sync (clone de dev en dev, ce clone-ci en prod).
WorkingDirectory=/home/robocar/robocar-Paris-PGE_MSC
Environment=OPENBLAS_CORETYPE=ARMV8
ExecStart=/usr/bin/python3 -u -m src.cam.hub --width 640 --height 320 --depth
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
